# var/lib/dpkg/info/fwupd.list
/.
/etc
/etc/fwupd.conf
/etc/pki
/etc/pki/fwupd-metadata
/etc/pki/fwupd-metadata/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd
/etc/pki/fwupd/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd/GPG-KEY-Hughski-Limited
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.fwupd.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fwupd-offline-update.service
/lib/systemd/system/system-update.target.wants
/lib/systemd/system/fwupd.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-fwupd-devices.rules
/usr
/usr/bin
/usr/bin/dfu-tool
/usr/bin/fwupdmgr
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_raspberrypi.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_udev.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_uefi.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_dfu.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_ebitdo.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_usb.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_steelseries.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_upower.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_colorhug.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_test.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_altos.so
/usr/lib/x86_64-linux-gnu/fwupd-plugins-2/libfu_plugin_unifying.so
/usr/lib/x86_64-linux-gnu/fwupd
/usr/lib/x86_64-linux-gnu/fwupd/fwupd
/usr/share
/usr/share/doc
/usr/share/doc/fwupd
/usr/share/doc/fwupd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dfu-tool.1.gz
/usr/share/man/man1/fwupdmgr.1.gz
/usr/share/polkit-1
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.fwupd.rules
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.fwupd.policy
/usr/share/app-info
/usr/share/app-info/xmls
/usr/share/app-info/xmls/org.freedesktop.fwupd.xml
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.fwupd.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.fwupd.service
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/fwupd.pkla
/var/lib/fwupd
/var/cache
/var/cache/app-info
/var/cache/app-info/xmls
/lib/systemd/system/system-update.target.wants/fwupd-offline-update.service
/usr/share/doc/fwupd/changelog.Debian.gz
