# var/lib/dpkg/info/ganglia-monitor.conffiles
/etc/ganglia/gmond.conf
/etc/init.d/ganglia-monitor
/etc/init/ganglia-monitor.conf
