# var/lib/dpkg/info/ganglia-monitor.list
/.
/usr
/usr/bin
/usr/bin/gstat
/usr/bin/gmetric
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/gmond.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/gmetric.1.gz
/usr/share/man/man1/gmond.1.gz
/usr/share/man/man1/gstat.1.gz
/usr/share/doc
/usr/share/doc/ganglia-monitor
/usr/share/doc/ganglia-monitor/README.Debian
/usr/share/doc/ganglia-monitor/copyright
/usr/share/doc/ganglia-monitor/NEWS.gz
/usr/sbin
/usr/sbin/gmond
/etc
/etc/init.d
/etc/init.d/ganglia-monitor
/etc/init
/etc/init/ganglia-monitor.conf
/etc/ganglia
/etc/ganglia/gmond.conf
/usr/share/doc/ganglia-monitor/changelog.Debian.gz
