# var/lib/dpkg/info/gridengine-exec.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/gridengine-exec
/usr/share/doc/gridengine-exec/NEWS.Debian.gz
/usr/share/doc/gridengine-exec/copyright
/usr/share/doc/gridengine-exec/changelog.Debian.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/sge_shepherd.8.gz
/usr/share/man/man8/sge_execd.8.gz
/usr/lib
/usr/lib/gridengine
/usr/lib/gridengine/sge_coshepherd
/usr/lib/gridengine/sge_execd
/usr/lib/gridengine/sge_shepherd
/usr/lib/gridengine/qrsh_starter
/usr/sbin
/etc
/etc/init.d
/etc/init.d/gridengine-exec
/usr/sbin/sge_coshepherd
/usr/sbin/sge_execd
/usr/sbin/sge_shepherd
