# var/lib/dpkg/info/kerneloops-daemon.list
/.
/usr
/usr/bin
/usr/bin/kerneloops-submit
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/kerneloops.8.gz
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_kerneloops.py
/usr/share/doc
/usr/share/doc/kerneloops-daemon
/usr/share/doc/kerneloops-daemon/copyright
/usr/share/doc/kerneloops-daemon/changelog.Debian.gz
/usr/sbin
/usr/sbin/kerneloops
/etc
/etc/default
/etc/default/kerneloops
/etc/kerneloops.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/kerneloops.dbus
/etc/init.d
/etc/init.d/kerneloops
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/kerneloops.service
