# var/lib/dpkg/info/nfs-common.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nfs-idmapd.service
/lib/systemd/system/rpc-svcgssd.service
/lib/systemd/system/auth-rpcgss-module.service
/lib/systemd/system/rpc-gssd.service
/lib/systemd/system/nfs-config.service
/lib/systemd/system/run-rpc_pipefs.mount
/lib/systemd/system/rpc-statd-notify.service
/lib/systemd/system/rpc-statd.service
/lib/systemd/system/nfs-client.target
/lib/systemd/system/proc-fs-nfsd.mount
/lib/systemd/system/nfs-utils.service
/etc
/etc/init
/etc/init/idmapd-mounting.conf
/etc/init/statd.override
/etc/init/statd.conf
/etc/init/gssd-mounting.conf
/etc/init/idmapd.conf
/etc/init/statd-mounting.override
/etc/init/statd-mounting.conf
/etc/init/gssd.conf
/etc/request-key.d
/etc/request-key.d/id_resolver.conf
/etc/init.d
/var
/var/lib
/var/lib/nfs
/var/lib/nfs/state
/var/lib/nfs/sm.bak
/var/lib/nfs/sm
/sbin
/sbin/showmount
/sbin/osd_login
/sbin/rpc.statd
/sbin/sm-notify
/sbin/mount.nfs
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/scripts
/usr/lib/systemd/scripts/nfs-utils_env.sh
/usr/sbin
/usr/sbin/mountstats
/usr/sbin/gss_clnt_send_err
/usr/sbin/rpc.gssd
/usr/sbin/nfsstat
/usr/sbin/rpc.svcgssd
/usr/sbin/rpcdebug
/usr/sbin/nfsidmap
/usr/sbin/start-statd
/usr/sbin/rpc.idmapd
/usr/sbin/blkmapd
/usr/sbin/gss_destroy_creds
/usr/sbin/nfsiostat
/usr/share
/usr/share/bug
/usr/share/bug/nfs-utils
/usr/share/bug/nfs-utils/control
/usr/share/bug/nfs-utils/presubj
/usr/share/bug/nfs-common
/usr/share/bug/nfs-common/control
/usr/share/bug/nfs-common/script
/usr/share/doc
/usr/share/doc/nfs-common
/usr/share/doc/nfs-common/copyright
/usr/share/doc/nfs-common/changelog.Debian.gz
/usr/share/doc/nfs-common/README.Debian.nfsv4
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/nfs-common
/usr/share/nfs-common
/usr/share/nfs-common/conffiles
/usr/share/nfs-common/conffiles/idmapd.conf
/usr/share/nfs-common/conffiles/nfs-common.default
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/nfsmount.conf.5.gz
/usr/share/man/man5/nfs.5.gz
/usr/share/man/man5/nfsidmap.5.gz
/usr/share/man/man8
/usr/share/man/man8/showmount.8.gz
/usr/share/man/man8/idmapd.8.gz
/usr/share/man/man8/blkmapd.8.gz
/usr/share/man/man8/nfsstat.8.gz
/usr/share/man/man8/svcgssd.8.gz
/usr/share/man/man8/statd.8.gz
/usr/share/man/man8/mount.nfs.8.gz
/usr/share/man/man8/nfsiostat.8.gz
/usr/share/man/man8/mountstats.8.gz
/usr/share/man/man8/rpcdebug.8.gz
/usr/share/man/man8/sm-notify.8.gz
/usr/share/man/man8/gssd.8.gz
/usr/share/man/man8/umount.nfs.8.gz
/lib/systemd/system/idmapd.service
/lib/systemd/system/gssd.service
/sbin/mount.nfs4
/sbin/umount.nfs
/sbin/umount.nfs4
/usr/share/man/man8/rpc.gssd.8.gz
/usr/share/man/man8/rpc.statd.8.gz
/usr/share/man/man8/rpc.svcgssd.8.gz
/usr/share/man/man8/rpc.idmapd.8.gz
/usr/share/man/man8/rpc.sm-notify.8.gz
