# var/lib/dpkg/info/nis.list
/.
/usr
/usr/share
/usr/share/locale
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/ypbind-mt.mo
/usr/share/locale/de/LC_MESSAGES/yp-tools.mo
/usr/share/nis
/usr/share/nis/nis.default
/usr/share/bug
/usr/share/bug/nis
/usr/share/bug/nis/script
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/nis
/usr/share/doc
/usr/share/doc/nis
/usr/share/doc/nis/examples
/usr/share/doc/nis/examples/ypserv.conf
/usr/share/doc/nis/nis.debian.howto.gz
/usr/share/doc/nis/copyright
/usr/share/doc/nis/changelog.Debian.gz
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/ypserv.conf.5.gz
/usr/share/man/man5/nicknames.5.gz
/usr/share/man/man5/yp.conf.5.gz
/usr/share/man/man5/netgroup.5.gz
/usr/share/man/man1
/usr/share/man/man1/ypchfn.1.gz
/usr/share/man/man1/ypwhich.1.gz
/usr/share/man/man1/ypmatch.1.gz
/usr/share/man/man1/ypchsh.1.gz
/usr/share/man/man1/ypcat.1.gz
/usr/share/man/man1/yppasswd.1.gz
/usr/share/man/man8
/usr/share/man/man8/ypdomainname.8.gz
/usr/share/man/man8/ypserv.8.gz
/usr/share/man/man8/yphelper.8.gz
/usr/share/man/man8/makedbm.8.gz
/usr/share/man/man8/rpc.yppasswdd.8.gz
/usr/share/man/man8/ypbind.8.gz
/usr/share/man/man8/domainname.8.gz
/usr/share/man/man8/ypxfrd.8.gz
/usr/share/man/man8/ypset.8.gz
/usr/share/man/man8/yppush.8.gz
/usr/share/man/man8/yppoll.8.gz
/usr/share/man/man8/ypinit.8.gz
/usr/share/man/man8/mknetid.8.gz
/usr/share/man/man8/pwupdate.8.gz
/usr/share/man/man8/nisdomainname.8.gz
/usr/share/man/man8/revnetgroup.8.gz
/usr/share/man/man8/rpc.ypxfrd.8.gz
/usr/share/man/man8/ypxfr.8.gz
/usr/share/man/man8/yppasswdd.8.gz
/usr/share/man/man8/yptest.8.gz
/usr/bin
/usr/bin/ypmatch
/usr/bin/ypcat
/usr/bin/ypwhich
/usr/bin/yppasswd
/usr/lib
/usr/lib/yp
/usr/lib/yp/yphelper
/usr/lib/yp/ypxfr_1perhour
/usr/lib/yp/ypxfr_2perday
/usr/lib/yp/pwupdate
/usr/lib/yp/ypinit
/usr/lib/yp/match_printcap
/usr/lib/yp/create_printcap
/usr/lib/yp/convert-ypserv-conf
/usr/lib/yp/ypxfr
/usr/lib/yp/mknetid
/usr/lib/yp/makedbm
/usr/lib/yp/ypxfr_1perday
/usr/lib/yp/revnetgroup
/usr/sbin
/usr/sbin/ypserv
/usr/sbin/rpc.yppasswdd
/usr/sbin/rpc.ypxfrd
/usr/sbin/ypbind
/usr/sbin/ypset
/usr/sbin/yppoll
/usr/sbin/ypserv_test
/usr/sbin/yppush
/usr/sbin/yptest
/bin
/etc
/etc/netgroup
/etc/yp.conf
/etc/ypserv.securenets
/etc/init.d
/etc/init.d/nis
/etc/ypserv.conf
/var
/var/yp
/var/yp/nicknames
/var/yp/Makefile
/usr/share/doc/nis/examples/nis.default
/usr/bin/ypchsh
/usr/bin/ypchfn
