# var/lib/dpkg/info/nis.conffiles
/etc/netgroup
/etc/init.d/nis
/etc/ypserv.conf
/etc/ypserv.securenets
/etc/yp.conf
/var/yp/Makefile
/var/yp/nicknames
