# var/lib/dpkg/info/nullmailer.conffiles
/etc/init.d/nullmailer
/etc/logcheck/ignore.d.paranoid/nullmailer
/etc/logcheck/ignore.d.server/nullmailer
/etc/logcheck/ignore.d.workstation/nullmailer
