# var/lib/dpkg/info/nullmailer.list
/.
/etc
/etc/logcheck
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/nullmailer
/etc/logcheck/ignore.d.workstation
/etc/logcheck/ignore.d.workstation/nullmailer
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/nullmailer
/etc/init.d
/etc/init.d/nullmailer
/etc/nullmailer
/var
/var/spool
/var/spool/nullmailer
/var/spool/nullmailer/tmp
/var/spool/nullmailer/queue
/usr
/usr/sbin
/usr/sbin/nullmailer-send
/usr/sbin/nullmailer-queue
/usr/sbin/sendmail
/usr/bin
/usr/bin/nullmailer-inject
/usr/bin/nullmailer-smtpd
/usr/bin/mailq
/usr/lib
/usr/lib/nullmailer
/usr/lib/nullmailer/qmqp
/usr/lib/nullmailer/smtp
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/nullmailer
/usr/share/doc
/usr/share/doc/nullmailer
/usr/share/doc/nullmailer/copyright
/usr/share/doc/nullmailer/changelog.Debian.gz
/usr/share/doc/nullmailer/TODO
/usr/share/doc/nullmailer/README.Debian
/usr/share/doc/nullmailer/README
/usr/share/doc/nullmailer/BUGS
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/nullmailer.7.gz
/usr/share/man/man8
/usr/share/man/man8/nullmailer-queue.8.gz
/usr/share/man/man8/nullmailer-smtpd.8.gz
/usr/share/man/man8/nullmailer-send.8.gz
/usr/share/man/man8/nullmailer-smtp.8.gz
/usr/share/man/man8/nullmailer-qmqp.8.gz
/usr/share/man/man1
/usr/share/man/man1/newaliases.1.gz
/usr/share/man/man1/nullmailer-inject.1.gz
/usr/share/man/man1/sendmail.1.gz
/usr/share/man/man5
/usr/bin/newaliases
/usr/lib/sendmail
/usr/share/doc/nullmailer/NEWS.gz
/usr/share/man/man1/mailq.1.gz
/usr/share/man/man5/pausetime.5.gz
/usr/share/man/man5/remotes.5.gz
/usr/share/man/man5/defaulthost.5.gz
/usr/share/man/man5/adminaddr.5.gz
/usr/share/man/man5/defaultdomain.5.gz
/usr/share/man/man5/idhost.5.gz
