# var/lib/dpkg/info/nvidia-prime.list
/.
/sbin
/sbin/prime-switch
/sbin/get-quirk-options
/sbin/prime-offload
/etc
/etc/init
/etc/init/nvidia-prime.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nvidia-prime.service
/usr
/usr/bin
/usr/bin/prime-supported
/usr/bin/prime-select
/usr/lib
/usr/lib/ubiquity
/usr/lib/ubiquity/dm-scripts
/usr/lib/ubiquity/dm-scripts/oem
/usr/lib/ubiquity/dm-scripts/oem/00-prime-offload
/usr/share
/usr/share/lightdm
/usr/share/lightdm/lightdm.conf.d
/usr/share/lightdm/lightdm.conf.d/90-nvidia.conf
/usr/share/nvidia-prime
/usr/share/nvidia-prime/prime-quirks
/usr/share/doc
/usr/share/doc/nvidia-prime
/usr/share/doc/nvidia-prime/changelog.gz
/usr/share/doc/nvidia-prime/copyright
