# var/lib/dpkg/info/pppconfig.list
/.
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/pppconfig
/usr/share/menu
/usr/share/menu/pppconfig
/usr/share/doc
/usr/share/doc/pppconfig
/usr/share/doc/pppconfig/changelog.gz
/usr/share/doc/pppconfig/copyright
/usr/share/man
/usr/share/man/de
/usr/share/man/de/man8
/usr/share/man/de/man8/pppconfig.8.gz
/usr/share/man/man8
/usr/share/man/man8/pppconfig.8.gz
/usr/share/man/pt
/usr/share/man/pt/man8
/usr/share/man/pt/man8/pppconfig.8.gz
/usr/share/man/fr
/usr/share/man/fr/man8
/usr/share/man/fr/man8/pppconfig.8.gz
/usr/sbin
/usr/sbin/pppconfig
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dns-clean.service
/etc
/etc/ppp
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/0dns-up
/etc/ppp/resolv
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/0dns-down
/etc/init.d
/etc/init.d/dns-clean
