# var/lib/dpkg/info/rpcbind.conffiles
/etc/init.d/rpcbind
/etc/init/portmap-wait.conf
/etc/init/rpcbind-boot.conf
/etc/init/rpcbind.conf
/etc/insserv.conf.d/rpcbind
