# var/lib/dpkg/info/rpcbind.list
/.
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/rpcbind
/usr/share/doc
/usr/share/doc/rpcbind
/usr/share/doc/rpcbind/copyright
/usr/share/doc/rpcbind/changelog.Debian.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/rpcinfo.8.gz
/usr/share/man/man8/rpcbind.8.gz
/usr/share/man/man7
/usr/share/man/man7/rpcinfo.7.gz
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/rpcbind.conf
/usr/sbin
/usr/sbin/rpcinfo
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rpcbind.socket
/lib/systemd/system/rpcbind.service
/sbin
/sbin/rpcbind
/etc
/etc/insserv.conf.d
/etc/insserv.conf.d/rpcbind
/etc/init.d
/etc/init.d/rpcbind
/etc/init
/etc/init/portmap-wait.conf
/etc/init/rpcbind.conf
/etc/init/rpcbind-boot.conf
/lib/systemd/system/portmap.service
