# var/lib/dpkg/info/sddm.conffiles
usr/share/sddm/scripts/Xsetup
usr/share/sddm/scripts/Xsession
usr/share/sddm/scripts/Xstop
usr/share/sddm/scripts/wayland-session
/etc/dbus-1/system.d/sddm_org.freedesktop.DisplayManager.conf
/etc/init.d/sddm
/etc/init/sddm.conf
/etc/pam.d/sddm
/etc/pam.d/sddm-autologin
/etc/pam.d/sddm-greeter
/etc/sddm/Xsession
