# var/lib/dpkg/info/xdiagnose.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/failsafe-graphical.target
/lib/systemd/system/display-manager.service.d
/lib/systemd/system/display-manager.service.d/xdiagnose.conf
/lib/systemd/system/failsafe-x.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-xdiagnose.rules
/usr
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/xdiagnose-3.8.4.1.egg-info
/usr/lib/python3/dist-packages/xdiagnose
/usr/lib/python3/dist-packages/xdiagnose/welcome.py
/usr/lib/python3/dist-packages/xdiagnose/x_pkg_names.py
/usr/lib/python3/dist-packages/xdiagnose/pci_devices.py
/usr/lib/python3/dist-packages/xdiagnose/application.py
/usr/lib/python3/dist-packages/xdiagnose/config_update.py
/usr/lib/python3/dist-packages/xdiagnose/edid.py
/usr/lib/python3/dist-packages/xdiagnose/xorglog.py
/usr/lib/python3/dist-packages/xdiagnose/utils
/usr/lib/python3/dist-packages/xdiagnose/utils/processes.py
/usr/lib/python3/dist-packages/xdiagnose/utils/url_io.py
/usr/lib/python3/dist-packages/xdiagnose/utils/math.py
/usr/lib/python3/dist-packages/xdiagnose/utils/debug.py
/usr/lib/python3/dist-packages/xdiagnose/utils/config.py
/usr/lib/python3/dist-packages/xdiagnose/utils/text.py
/usr/lib/python3/dist-packages/xdiagnose/utils/screen.py
/usr/lib/python3/dist-packages/xdiagnose/utils/lists.py
/usr/lib/python3/dist-packages/xdiagnose/utils/node.py
/usr/lib/python3/dist-packages/xdiagnose/utils/readurl.py
/usr/lib/python3/dist-packages/xdiagnose/utils/dicts.py
/usr/lib/python3/dist-packages/xdiagnose/utils/paths.py
/usr/lib/python3/dist-packages/xdiagnose/utils/__init__.py
/usr/lib/python3/dist-packages/xdiagnose/utils/json_io.py
/usr/lib/python3/dist-packages/xdiagnose/utils/option_handler.py
/usr/lib/python3/dist-packages/xdiagnose/utils/dates.py
/usr/lib/python3/dist-packages/xdiagnose/utils/file_io.py
/usr/lib/python3/dist-packages/xdiagnose/errors_treeview.py
/usr/lib/python3/dist-packages/xdiagnose/applet.py
/usr/lib/python3/dist-packages/xdiagnose/__init__.py
/usr/lib/python3/dist-packages/xdiagnose/diagnostics.py
/usr/lib/python3/dist-packages/xdiagnose/info.py
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_xorg.py
/usr/share/apport/apport-gpu-error-intel.py
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/xdiagnose.svg
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/com.canonical.xdiagnose.policy
/usr/share/doc
/usr/share/doc/xdiagnose
/usr/share/doc/xdiagnose/copyright
/usr/share/doc/xdiagnose/README
/usr/share/doc/xdiagnose/changelog.gz
/usr/share/xdiagnose
/usr/share/xdiagnose/failsafeXinit
/usr/share/xdiagnose/failsafeXServer
/usr/share/xdiagnose/workloads
/usr/share/xdiagnose/workloads/do_screensaver_loop
/usr/share/xdiagnose/workloads/run_workloads
/usr/share/xdiagnose/workloads/do_kernel_compile_loop
/usr/share/xdiagnose/workloads/repro.sh
/usr/share/xdiagnose/workloads/do_cpu_spin_loop
/usr/share/xdiagnose/workloads/do_monitor_disable_loop
/usr/share/xdiagnose/workloads/do_chws_loop
/usr/share/xdiagnose/workloads/do_monitor_rotation_loop
/usr/share/xdiagnose/workloads/youtube-reload.html
/usr/share/xdiagnose/workloads/do_monitor_resolution_loop
/usr/share/xdiagnose/workloads/README
/usr/share/xdiagnose/workloads/do_disk_write_loop
/usr/share/xdiagnose/workloads/do_vtswitch_loop
/usr/share/xdiagnose/workloads/do_glx_loop
/usr/share/xdiagnose/workloads/do_video_loop
/usr/share/xdiagnose/workloads/youtube-loop.html
/usr/share/xdiagnose/bubu.lib
/usr/share/applications
/usr/share/applications/xdiagnose.desktop
/usr/bin
/usr/bin/xrandr-tool
/usr/bin/xrotate
/usr/bin/xdiagnose-pkexec
/usr/bin/dpkg-log-summary
/usr/bin/xpci
/usr/bin/xdiagnose
/usr/bin/xedid
/etc
/etc/init
/etc/init/failsafe-x.conf
/usr/share/apport/package-hooks/source_xserver-xorg-video-sisusb.py
/usr/share/apport/package-hooks/source_x11proto-xinerama.py
/usr/share/apport/package-hooks/source_libdmx.py
/usr/share/apport/package-hooks/source_libxss.py
/usr/share/apport/package-hooks/xserver-xorg-video-vesa-lts-quantal.py
/usr/share/apport/package-hooks/xserver-xorg-input-synaptics-lts-quantal.py
/usr/share/apport/package-hooks/source_libxxf86vm.py
/usr/share/apport/package-hooks/source_x11proto-fixes.py
/usr/share/apport/package-hooks/source_x11proto-render.py
/usr/share/apport/package-hooks/source_xauth.py
/usr/share/apport/package-hooks/source_libxvmc.py
/usr/share/apport/package-hooks/source_x11proto-randr.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-keyboard.py
/usr/share/apport/package-hooks/source_xcompmgr.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-vmmouse.py
/usr/share/apport/package-hooks/source_x11-xfs-utils.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-s3.py
/usr/share/apport/package-hooks/xserver-xorg-video-neomagic-lts-quantal.py
/usr/share/apport/package-hooks/source_libxdamage.py
/usr/share/apport/package-hooks/xserver-xorg-video-sisusb-lts-quantal.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-displaylink.py
/usr/share/apport/package-hooks/source_libwacom.py
/usr/share/apport/package-hooks/source_x11proto-xf86bigfont.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-dummy.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-fbdev.py
/usr/share/apport/package-hooks/source_fglrx-installer-updates.py
/usr/share/apport/package-hooks/source_x11-session-utils.py
/usr/share/apport/package-hooks/xserver-xorg-video-cirrus-lts-quantal.py
/usr/share/apport/package-hooks/source_xf86-input-evtouch.py
/usr/share/apport/package-hooks/xserver-xorg-video-modesetting-lts-quantal.py
/usr/share/apport/package-hooks/xserver-xorg-video-vmware-lts-quantal.py
/usr/share/apport/package-hooks/xorg-server-lts-quantal.py
/usr/share/apport/package-hooks/source_x11proto-print.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-evdev.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-qxl.py
/usr/share/apport/package-hooks/source_libxrandr.py
/usr/share/apport/package-hooks/source_x-kit.py
/usr/share/apport/package-hooks/source_x11-xserver-utils.py
/usr/share/apport/package-hooks/mesa-lts-quantal.py
/usr/share/apport/package-hooks/source_libfontenc.py
/usr/share/apport/package-hooks/source_xft.py
/usr/share/apport/package-hooks/source_xf86-video-omapfb.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-ati.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-synaptics.py
/usr/share/apport/package-hooks/source_xtrans.py
/usr/share/apport/package-hooks/source_libpciaccess.py
/usr/share/apport/package-hooks/source_libfs.py
/usr/share/apport/package-hooks/source_libxi.py
/usr/share/apport/package-hooks/source_x11proto-xf86dga.py
/usr/share/apport/package-hooks/source_libxcb.py
/usr/share/apport/package-hooks/source_libxres.py
/usr/share/apport/package-hooks/xserver-xorg-video-savage-lts-quantal.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-amdgpu.py
/usr/share/apport/package-hooks/xserver-xorg-video-nouveau-lts-quantal.py
/usr/share/apport/package-hooks/xserver-xorg-video-mga-lts-quantal.py
/usr/share/apport/package-hooks/source_xinput-calibrator.py
/usr/share/apport/package-hooks/source_xorg-options-editor-gtk.py
/usr/share/apport/package-hooks/source_xfs.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-psb.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-joystick.py
/usr/share/apport/package-hooks/source_x11proto-gl.py
/usr/share/apport/package-hooks/xserver-xorg-video-tdfx-lts-quantal.py
/usr/share/apport/package-hooks/source_x11-utils.py
/usr/share/apport/package-hooks/source_radeontool.py
/usr/share/apport/package-hooks/libxrandr-lts-quantal.py
/usr/share/apport/package-hooks/source_x11proto-record.py
/usr/share/apport/package-hooks/source_nvidia-graphics-drivers-352.py
/usr/share/apport/package-hooks/source_libxfixes.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-nouveau.py
/usr/share/apport/package-hooks/xf86-input-mtrack-lts-quantal.py
/usr/share/apport/package-hooks/source_x11-xkb-utils.py
/usr/share/apport/package-hooks/xserver-xorg-input-joystick-lts-quantal.py
/usr/share/apport/package-hooks/source_libx11.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-nv.py
/usr/share/apport/package-hooks/source_xfonts-utils.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-void.py
/usr/share/apport/package-hooks/source_linux-restricted-modules-2.6.24.py
/usr/share/apport/package-hooks/source_xdm.py
/usr/share/apport/package-hooks/source_xinput.py
/usr/share/apport/package-hooks/source_x11proto-xcmisc.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-sis.py
/usr/share/apport/package-hooks/source_libxfont.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-mouse.py
/usr/share/apport/package-hooks/source_x11proto-xf86dri.py
/usr/share/apport/package-hooks/source_x11-apps.py
/usr/share/apport/package-hooks/source_x11proto-fontcache.py
/usr/share/apport/package-hooks/source_xrandr.py
/usr/share/apport/package-hooks/source_libice.py
/usr/share/apport/package-hooks/source_x11proto-trap.py
/usr/share/apport/package-hooks/xserver-xorg-video-ati-lts-quantal.py
/usr/share/apport/package-hooks/source_nvidia-graphics-drivers-304.py
/usr/share/apport/package-hooks/xserver-xorg-input-void-lts-quantal.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-vesa.py
/usr/share/apport/package-hooks/source_wayland.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-libinput.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-intel.py
/usr/share/apport/package-hooks/source_xrestop.py
/usr/share/apport/package-hooks/source_x11proto-scrnsaver.py
/usr/share/apport/package-hooks/source_libxft.py
/usr/share/apport/package-hooks/source_libxrender.py
/usr/share/apport/package-hooks/source_xorg-server.py
/usr/share/apport/package-hooks/source_libvdpau.py
/usr/share/apport/package-hooks/source_nvidia-graphics-drivers-352-updates.py
/usr/share/apport/package-hooks/source_x11proto-dmx.py
/usr/share/apport/package-hooks/source_xfonts-core.py
/usr/share/apport/package-hooks/source_xtrace.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-geode.py
/usr/share/apport/package-hooks/source_nvidia-graphics-drivers-340-updates.py
/usr/share/apport/package-hooks/source_weston.py
/usr/share/apport/package-hooks/xserver-xorg-video-geode-lts-quantal.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-vmware.py
/usr/share/apport/package-hooks/source_libxaw.py
/usr/share/apport/package-hooks/source_x11proto-bigreqs.py
/usr/share/apport/package-hooks/xserver-xorg-core.py
/usr/share/apport/package-hooks/source_x11proto-input.py
/usr/share/apport/package-hooks/source_xcb-proto.py
/usr/share/apport/package-hooks/source_libxdmcp.py
/usr/share/apport/package-hooks/xserver-xorg-input-mouse-lts-quantal.py
/usr/share/apport/package-hooks/source_libxkbcommon.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-siliconmotion.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-savage.py
/usr/share/apport/package-hooks/xserver-xorg-video-r128-lts-quantal.py
/usr/share/apport/package-hooks/source_pixman.py
/usr/share/apport/package-hooks/source_libxxf86dga.py
/usr/share/apport/package-hooks/source_libsm.py
/usr/share/apport/package-hooks/xserver-xorg-video-mach64-lts-quantal.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-r128.py
/usr/share/apport/package-hooks/source_libxv.py
/usr/share/apport/package-hooks/source_x11proto-core.py
/usr/share/apport/package-hooks/source_libxau.py
/usr/share/apport/package-hooks/source_wacom-tools.py
/usr/share/apport/package-hooks/xserver-xorg-video-trident-lts-quantal.py
/usr/share/apport/package-hooks/source_xorg-sgml-doctools.py
/usr/share/apport/package-hooks/source_nvidia-graphics-drivers-340.py
/usr/share/apport/package-hooks/source_xbitmaps.py
/usr/share/apport/package-hooks/x11-xserver-utils-lts-quantal.py
/usr/share/apport/package-hooks/source_fglrx-installer-experimental-9.py
/usr/share/apport/package-hooks/xserver-xorg-video-dummy-lts-quantal.py
/usr/share/apport/package-hooks/xserver-xorg-video-openchrome-lts-quantal.py
/usr/share/apport/package-hooks/source_fglrx-installer.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-trident.py
/usr/share/apport/package-hooks/source_libxkbui.py
/usr/share/apport/package-hooks/source_xf86-input-wacom.py
/usr/share/apport/package-hooks/source_libxcursor.py
/usr/share/apport/package-hooks/source_mesa.py
/usr/share/apport/package-hooks/xserver-xorg-video-intel-lts-quantal.py
/usr/share/apport/package-hooks/source_xorg-docs.py
/usr/share/apport/package-hooks/source_x11proto-resource.py
/usr/share/apport/package-hooks/xserver-xorg-input-vmmouse-lts-quantal.py
/usr/share/apport/package-hooks/source_xterm.py
/usr/share/apport/package-hooks/source_x11proto-xext.py
/usr/share/apport/package-hooks/source_xutils-dev.py
/usr/share/apport/package-hooks/source_x11proto-dri2.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-openchrome.py
/usr/share/apport/package-hooks/source_x11proto-xf86vidmode.py
/usr/share/apport/package-hooks/source_x11proto-video.py
/usr/share/apport/package-hooks/source_libxkbfile.py
/usr/share/apport/package-hooks/source_x11proto-kb.py
/usr/share/apport/package-hooks/xserver-xorg-video-sis-lts-quantal.py
/usr/share/apport/package-hooks/xserver-xorg-video-siliconmotion-lts-quantal.py
/usr/share/apport/package-hooks/source_libxtst.py
/usr/share/apport/package-hooks/source_nvidia-graphics-drivers-304-updates.py
/usr/share/apport/package-hooks/source_libdrm.py
/usr/share/apport/package-hooks/xserver-xorg-video-s3-lts-quantal.py
/usr/share/apport/package-hooks/source_x11proto-damage.py
/usr/share/apport/package-hooks/source_x11proto-fonts.py
/usr/share/apport/package-hooks/source_libxt.py
/usr/share/apport/package-hooks/source_libxcomposite.py
/usr/share/apport/package-hooks/source_xf86-video-displaylink.py
/usr/share/apport/package-hooks/source_libxmu.py
/usr/share/apport/package-hooks/source_mesa-demos.py
/usr/share/apport/package-hooks/source_xserver-xorg-input-wacom.py
/usr/share/apport/package-hooks/source_xf86-input-mtrack.py
/usr/share/apport/package-hooks/xorg-lts-quantal.py
/usr/share/apport/package-hooks/source_x11proto-composite.py
/usr/share/apport/package-hooks/source_libxfontcache.py
/usr/share/apport/package-hooks/xf86-input-wacom-lts-quantal.py
/usr/share/apport/package-hooks/source_xf86-video-msm.py
/usr/share/apport/package-hooks/xserver-xorg-input-keyboard-lts-quantal.py
/usr/share/apport/package-hooks/source_x11proto-evie.py
/usr/share/apport/package-hooks/source_xinit.py
/usr/share/apport/package-hooks/source_libxext.py
/usr/share/apport/package-hooks/source_libxinerama.py
/usr/share/apport/package-hooks/source_xkeyboard-config.py
/usr/share/apport/package-hooks/source_xcursor-themes.py
/usr/share/apport/package-hooks/source_xserver-xorg-video-mga.py
/usr/share/apport/package-hooks/source_libxpm.py
